# launchd/com.alexrangelalves.pausasbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexrangelalves.pausasbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>/Users/alexrangelalves/Downloads/Conexão mtds/pausas_bot/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/alexrangelalves/Downloads/Conexão mtds/pausas_bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/alexrangelalves/Downloads/Conexão mtds/pausas_bot/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alexrangelalves/Downloads/Conexão mtds/pausas_bot/bot.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Library/Frameworks/Python.framework/Versions/3.14/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
